# backend/platform/docker/ansible/roles/supervisor/templates/project_supervisor.conf.j2
[program:django]
command=gunicorn config.wsgi -c /etc/gunicorn.conf.py
{% include "cmd_vars.j2" %}

[group:{{ project_name }}]
programs=django
